# scheduling/daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.job-search-agent.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>PLACEHOLDER_SCRIPT_PATH</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>PLACEHOLDER_HOUR</integer>
        <key>Minute</key>
        <integer>PLACEHOLDER_MINUTE</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/job-search-agent-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/job-search-agent-launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
