# launchd/com.jaron.job-search-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jaron.job-search-health</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/bin/health-check.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSIDIAN_VAULT</key>
    <string>__OBSIDIAN_VAULT__</string>
    <key>REPO_ROOT</key>
    <string>__REPO_ROOT__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- Runs every hour, on the hour. -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>/tmp/job-search-agent-health.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/job-search-agent-health.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
