# launchd/com.jaron.job-search-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jaron.job-search-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/bin/run-weekly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSIDIAN_VAULT</key>
    <string>__OBSIDIAN_VAULT__</string>
    <key>REPO_ROOT</key>
    <string>__REPO_ROOT__</string>
    <key>CLAUDE_BIN</key>
    <string>__CLAUDE_BIN__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- Monday = Weekday 1 in launchd's StartCalendarInterval. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/job-search-agent-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/job-search-agent-launchd.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
