# launchd/com.jaron.job-search-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Template — values like __REPO_ROOT__, __OBSIDIAN_VAULT__, __CLAUDE_BIN__
     are substituted by bin/install.sh at install time. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jaron.job-search-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/bin/run-daily.sh</string>
  </array>

  <!-- Working directory is the REPO, NOT the vault.
       The vault is TCC-protected; using it as cwd causes launchd to fail
       with "getcwd: cannot access parent directories: Operation not permitted". -->
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSIDIAN_VAULT</key>
    <string>__OBSIDIAN_VAULT__</string>
    <key>REPO_ROOT</key>
    <string>__REPO_ROOT__</string>
    <key>CLAUDE_BIN</key>
    <string>__CLAUDE_BIN__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/job-search-agent-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/job-search-agent-launchd.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
